# src/Deploy/Resources/worker/supervisord.rootless.conf
; Canonical ROOTLESS supervisord config for the Vortos deploy-in-image (RideColor) model.
;
; In deploy-in-image one image serves every role and runs as the image's non-root user (e.g.
; `USER app`, uid 1000). A supervisord config with `user=root` — or a pidfile/socket/logfile under a
; root-owned directory (`/var/run`, `/var/log/supervisor`) — makes supervisord try to drop privileges
; as a non-root user and die: "Error: Can't drop privilege as nonroot user" (crash loop).
;
; This config is rootless by construction:
;   - no `[supervisord] user=` directive (supervisord runs as whoever launched it — the image user)
;   - pidfile + admin socket live in /tmp (writable by any uid)
;   - the main supervisord log goes to stdout (captured by `docker logs`), never a root-owned dir
;   - each program logs to stdout/stderr (12-factor) rather than files under /var/log
;
; Copy this into your image (e.g. `COPY supervisord.rootless.conf /etc/supervisord.conf`) or base your
; own on it. `RootlessWorkerCheck` (deploy:doctor) fails closed if the active worker config is rootful.

[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Example consumer program — replace the command with your real worker(s) ──
; Multiple [program:*] blocks (or a [program:*] with numprocs) spawn the queue/DLQ consumers. Each
; inherits the image's non-root user; do NOT add `user=` here.
[program:consumer]
command=php /app/bin/console messenger:consume async --time-limit=3600
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
